# files/supervisor-ivideon-server.conf
[program:ivideon-server]
command=ivideon-server
redirect_stdout=false
redirect_stderr=false
